# profile.d/android.sh
##
# setup android
##

ANDROID_SDK=/Developer/Android/android-sdk-mac_x86

PATH=${PATH}:${ANDROID_SDK}/tools:${ANDROID_SDK}/platform-tools

export ANDROID_SDK PATH
