# profile.d/brooklyn.sh
##
# Brooklyn configuration
##

export BROOKLYN_HOME=/opt/brooklyn
export BROOKLYN_CLASSPATH=/opt/brooklyn/conf

# vim:ts=4:sw=4
