# profile.d/environment.sh
##
# Environment
#
# Environment variables not exported elsewhere
#
# Copyright 2011-2015 by Andrew Donald Kennedy; All Rights Reserved.
##

# Directory where Atom.app is located
export ATOM_PATH="/Applications/Developer"

# vim:ts=4:sw=4
