# profile.d/go.sh
##
# Go environment variables
#
# Copyright 2016 by Andrew Donald Kennedy; All Rights Reserved.
##

export GOPATH=/Users/adk/Workspace/gocode
export PATH=${PATH}:${GOPATH}/bin

# vim:ts=4:sw=4
