# profile.d/limits.sh
#!/bin/bash
# set limits

# open file handles
ulimit -n 8192
