# profile.d/limits.sh
##
# Limits
#
# Copyright 2011-2015 by Andrew Donald Kennedy; All Rights Reserved.
##

# open file handles
ulimit -n 8192

# vim:ts=4:sw=4
