# profile.d/python.sh
##
# Python Configuration
#
# Copyright 2011-2017 by Andrew Donald Kennedy; All Rights Reserved.
##

##
# Python Environment
##
PYTHON27PATH="/opt/local/Library/Frameworks/Python.framework/Versions/2.7/lib/python2.7/site-packages"
PYTHON35PATH="/opt/local/Library/Frameworks/Python.framework/Versions/3.5/lib/python3.5/site-packages"
PYTHONPATH=${PYTHON35PATH}
#PYTHONVERBOSE=1
#PYTHONDEBUG=1
PYTHONCASEOK=true
PYTHONHTTPSVERIFY=0
export PYTHONPATH PYTHON27PATH PYTHON35PATH PYTHONVERBOSE PYTHONDEBUG PYTHONCASEOK PYTHONHTTPSVERIFY

# vim:ts=4:sw=4
