# profile.d/rust.sh
##
# Rust Configuration
#
# Copyright 2017 by Andrew Donald Kennedy; All Rights Reserved.
##

PATH=${PATH}:${HOME}/.cargo/bin
export PATH

# vim:ts=4:sw=4
