# profile.d/rvm.sh
##
# RVM
#
# Copyright 2011-2015 by Andrew Donald Kennedy; All Rights Reserved.
##

#set -x # debug

# Load RVM into a shell session *as a function*
[[ -s "$HOME/.rvm/scripts/rvm" ]] && source "$HOME/.rvm/scripts/rvm"

# vim:ts=4:sw=4
