# profile.d/vi.sh
##
# Set default editor to Vim
#
# Copyright 2011-2015 by Andrew Donald Kennedy; All Rights Reserved.
##

EDITOR=$(which vim)
VISUAL=$(which vim)

export EDITOR VISUAL

alias vi=vim

# vim:ts=4:sw=4
