# profile.d/vi.sh
##
# set default editor to vi
##

EDITOR=$(which vi)
VISUAL=$(which vi)

export EDITOR VISUAL
